# deploy/sommelier_frontend.service
[Unit]
Description=Sommelier Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=sommelier
WorkingDirectory=/srv/sommelier_frontend
ExecStart=/usr/bin/node /srv/sommelier_frontend/.next/standalone/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
